# ansible/roles/anubis/templates/anubis.service.j2
# =============================================================================
# Anubis WAF Systemd Service
# =============================================================================
# Purpose: Manage the Anubis service as a system daemon
# Generated by Ansible - DO NOT edit manually
#
# Usage:
#   systemctl start anubis
#   systemctl status anubis
#   journalctl -u anubis -f
# =============================================================================

[Unit]
Description=Anubis WAF - Proof-of-Work Challenge
Documentation=https://github.com/TecharoHQ/anubis
After=network.target

[Service]
Type=simple

# Unprivileged user for security
User=anubis
Group=anubis

# Execute binary with configuration
ExecStart=/usr/local/bin/anubis --config /etc/anubis/config.yml

# Auto-restart if the service crashes
Restart=always
RestartSec=5

# Memory limit for t4g.nano (512MB RAM + 1GB swap)
# 150M is sufficient for Anubis
MemoryMax=150M

# Security protections
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/etc/anubis

[Install]
WantedBy=multi-user.target
